# system/videoplayerd.service
[Unit]
Description=USB Video Autoplay Service
After=local-fs.target network.target
Wants=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/videoplayerd
WorkingDirectory=/usr/local/share/videoplayerd
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
